# pb2/crontab


@reboot sleep 60 && php /home/user/bin/send_ip.php
